# archiso/airootfs/etc/profile.d/cinnamon.sh
gsettings set org.cinnamon.desktop.background picture-uri 'file:///usr/share/backgrounds/gnome/truchet-l.jpg'
