# archiso/airootfs/etc/profile.d/cinnamon.sh
gsettings set org.cinnamon.desktop.background picture-options 'none'
gsettings set org.cinnamon.desktop.background color-shading-type 'solid'
gsettings set org.cinnamon.desktop.background primary-color '#2e2e34344040'
gsettings set org.cinnamon.theme name 'Nordic'
gsettings set org.cinnamon.desktop.wm.preferences theme 'Nordic'
gsettings set org.cinnamon.desktop.interface icon-theme 'Papirus-nordic-folders'
gsettings set org.cinnamon.desktop.interface gtk-theme 'Nordic'
gsettings set org.cinnamon.desktop.interface cursor-theme 'OpenZone_White'
gsettings set org.cinnamon.desktop.applications.terminal exec 'kitty'
